# source-code/DEBIAN/prerm
#!/bin/bash
set -e

# Stop all CasaOS services
SERVICES="casaos casaos-app-management casaos-local-storage casaos-user-service casaos-gateway casaos-message-bus"

for service in $SERVICES; do
    if command -v rc-service >/dev/null 2>&1; then
        rc-service "$service" stop 2>/dev/null || true
        rc-update del "$service" default 2>/dev/null || true
    else
        /etc/init.d/"$service" stop 2>/dev/null || true
    fi
done

exit 0
