# source-code/etc/init.d/casaos-local-storage
#!/sbin/openrc-run
# Copyright 2025 CasaOS Team
# Distributed under the terms of the GNU General Public License v2

name="CasaOS Local Storage Service"
description="CasaOS Local Storage - File management and storage"

command="/usr/bin/casaos-local-storage"
command_args="-c /etc/casaos/local-storage.conf"
command_background=true
pidfile="/var/run/casaos/local-storage.pid"

depend() {
	need net
	after casaos-message-bus casaos-gateway casaos-user-service
}

start_pre() {
	checkpath --directory --mode 0755 /var/run/casaos
}
